# systemd/leanframe-sync.service
[Unit]
Description=LeanFrame cloud sync (rclone)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/leanframe
User=leanframe
ExecStart=/opt/leanframe/.venv/bin/python -c "from photoframe.config import AppCfg; from photoframe.indexer import Library; from photoframe.sync import Syncer; from pathlib import Path; cfg=AppCfg.load(Path('config/leanframe.yaml')); lib=Library(cfg.paths.db, cfg.paths.library); Syncer(cfg, lib).run_all()"